# overlay/files/etc/periodic/15min/motd
#!/bin/sh
#. /etc/os-release

cat /etc/issue | sed 's@\\r@'"$(uname -r)"'@g' | sed 's@\\m@'"$(uname -m)"'@g' | sed 's@(\\l)@@g' > /etc/motd