# overlay/files/etc/systemd/system/getty@ttyS0.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/edgenode-dashboard
StandardInput=tty
StandardOutput=tty